# mariadb-columnstore/service/rsyslogd/run
#!/bin/sh
rm -f /var/run/syslogd.pid
exec /sbin/rsyslogd -n 2>&1
