# mariadb-columnstore/service/columnstore/finish
#!/bin/sh
rm -f /var/lock/subsys/columnstore
rm -f /var/lock/subsys/mysql-Columnstore
