# mariadb-columnstore/service/systemd-journald/run
#!/bin/sh
rm -rf /run/log/journal/*
/usr/lib/systemd/systemd-journald 2>&1
