# mariadb-columnstore/service/mysqld/run
#!/bin/sh
sleep 30
exec /usr/local/mariadb/columnstore/mysql/bin/mysqld --user=root
